# src/main/resources/bin/systemd/smart119-backend.service
[Unit]
Description=The smart119-backend system service
After=network.target

[Service]
LimitNOFILE=100000
Type=forking
ExecStart=/usr/bin/sh -c "/opt/smart119/smart119_backend/bin/smart119-backend.sh start"

[Install]
WantedBy=multi-user.target
